# udev/etc/udev/rules.d/98-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY}=":0", RUN+="/bin/bash $HOME/dotfiles/monitor.sh"
